# update-cached-status.cron
*/5 * * * * cd /home/access/public_html/cgi-bin && ruby update-cached-status.rb

